# systemd/strapi.service.template
[Unit]
Description=Strapi CMS for Edukacja Filozoficzna
After=meilisearch.service

[Service]
Type=simple
User=${USER}
Restart=on-failure
WorkingDirectory=${REPO}
ExecStart=${REPO}/scripts/run_jailed.sh strapi

[Install]
WantedBy=multi-user.target
